# backup/homelab.timer
[Unit]
Description=Run container backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
